# box/roles/oblac-drives-user/templates/etc/systemd/system/getty@tty1.service.d/override.conf.j2
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noissue --autologin {{ oblac_drives.user }} %I $TERM
Type=idle
