# box/roles/oblac-drives-update-service/templates/etc/systemd/system/docker.service.d/override.conf.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H unix:///var/run/docker.sock -H tcp://0.0.0.0:{{ box_docker_engine_api_port }} --api-cors-header="*"
