# box/roles/box-user/templates/etc/systemd/system/getty@tty1.service.d/override.conf.j2
[Service]
ExecStart=
ExecStart=-/sbin/agetty --noissue --autologin {{ box_user }} %I $TERM
Type=idle
